# systemdfiles/mongod.service
[Unit]
Description=MongoDB Database Service
After=network.target

[Service]
Type=forking
PIDFile=/var/run/mozdefdb/mozdefdb.pid
ExecStart=/usr/bin/mongod --config /etc/mongod.conf
ExecReload=/bin/kill -HUP $MAINPID
<<<<<<< HEAD
Restart=always
=======
#Restart=always
>>>>>>> 8d519538b9afa92609f5afa587256e02a25a554e
User=mozdef
Group=mozdef
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
